# config/systemd/puma.service.erb
[Unit]
Description=SonicD Puma
After=network.target

[Service]
Type=notify
WatchdogSec=10

User=<%= fetch(:user) %>
Environment=RAILS_ENV=<%= fetch(:stage)%>
Environment=ACTION_CABLE_ALLOWED_REQUEST_ORIGINS=<%= ENV['ACTION_CABLE_ALLOWED_REQUEST_ORIGINS'] %>
Environment=RAILS_LOG_TO_STDOUT=1
Environment=MALLOC_ARENA_MAX=2
WorkingDirectory=<%= current_path %>
ExecStart=<%= "/home/#{fetch(:user)}/.rvm/bin/rvm #{fetch(:rvm_ruby_version)} do #{current_path}/bin/bundle exec puma -C #{current_path}/config/puma.rb" %>
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always
KillMode=mixed

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sonicd-<%= fetch(:stage) %>-puma

[Install]
WantedBy=multi-user.target