# config/systemd/puma.service.erb
[Unit]
Description=SonicD Puma
After=network.target

[Service]
Type=simple
User=<%= fetch(:user) %>
Environment=RAILS_ENV=<%= fetch(:stage)%>
Environment=MALLOC_ARENA_MAX=2
WorkingDirectory=<%= current_path %>
ExecStart=<%= "/home/#{fetch(:user)}/.rvm/bin/rvm #{fetch(:rvm_ruby_version)} do #{current_path}/bin/bundle exec puma -C #{current_path}/config/puma.rb" %>
ExecStop=<%= "/home/#{fetch(:user)}/.rvm/bin/rvm #{fetch(:rvm_ruby_version)} do #{current_path}/bin/bundle exec pumactl -F #{current_path}/config/puma.rb stop" %>
ExecReload=<%= "/home/#{fetch(:user)}/.rvm/bin/rvm #{fetch(:rvm_ruby_version)} do #{current_path}/bin/bundle exec pumactl -F #{current_path}/config/puma.rb -S #{shared_path}/tmp/pids/puma.state phased-restart" %>
TimeoutSec=15
Restart=always
KillMode=process

[Install]
WantedBy=multi-user.target