# config/systemd/sidekiq.service.erb
[Unit]
Description=SonicD Sidekiq
After=network.target

[Service]
Type=notify
WatchdogSec=10

WorkingDirectory=<%= current_path %>
ExecStart=<%= "/home/#{fetch(:user)}/.rvm/bin/rvm #{fetch(:rvm_ruby_version)} do #{current_path}/bin/bundle exec sidekiq -e #{fetch(:stage)}" %>
User=<%= fetch(:user) %>
UMask=0002

# Greatly reduce Ruby memory fragmentation and heap usage
# https://www.mikeperham.com/2018/04/25/taming-rails-memory-bloat/
Environment=MALLOC_ARENA_MAX=2

RestartSec=1
Restart=on-failure

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sonicd-<%= fetch(:stage) %>-sidekiq

[Install]
WantedBy=multi-user.target